# system/files/private/etc/profile.d/sh.local##os.Darwin
../../../etc/profile.d/sh.local##os.Linux