# system/files/private/etc/profile.d/sitstart.sh##os.Darwin
../../../etc/profile.d/sitstart.sh##os.Linux